# polygonal_zones_editor/rootfs/etc/services.d/web/run
#!/usr/bin/with-contenv bashio
# Launch the addon web service as the unprivileged 'app' user. s6-overlay
# starts this after cont-init.d has chowned /data.

cd /app
exec s6-setuidgid app python main.py
