# polygonal_zones_editor/rootfs/etc/services.d/web/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# Launch the addon web service. Prefer dropping to the unprivileged 'app'
# user via s6-setuidgid; fall back to running as root with a loud warning
# if the wiring isn't available, so the service never fails to start for
# privilege-drop reasons alone.

# shellcheck disable=SC2164
# `cd /app` can technically fail, but `/app` is Dockerfile-installed and
# always present in a shipped image. If cd were to fail, the subsequent
# `python main.py` would also fail (main.py isn't at /), s6 would mark
# the service down, finish runs, and Supervisor restarts the addon —
# same end state as an explicit `|| exit`, with no live-Supervisor
# behaviour change to verify in this PR.
cd /app

if ! command -v s6-setuidgid >/dev/null 2>&1; then
    bashio::log.error "s6-setuidgid not found; running as root"
    exec python main.py
fi

if ! id app >/dev/null 2>&1; then
    bashio::log.error "'app' user not found; running as root"
    exec python main.py
fi

# Verify the app user can actually read options.json; if not, running as
# root is preferable to a startup crash.
if [ -f /data/options.json ] && ! s6-setuidgid app sh -c 'test -r /data/options.json'; then
    bashio::log.error "options.json not readable by 'app'; running as root"
    exec python main.py
fi

exec s6-setuidgid app python main.py
